# supervisord/supervisord.conf
[program:rpstomp]
command=/usr/bin/python /home/pi/raspberry-stomp/rpstomp.py